# docker/apache/supercronic/crontab
# Run every minute
*/1 * * * * php bin/console app:appointment:check
